# .docker/container/etc/supervisor/program/angular.conf
[program:angular]
; Angular dev server (DEV: ng serve) or SSR server (PROD: node server.mjs).
; Exits immediately with code 0 if DKZ_NODEJS_ANGULAR_VERSION_INSTALL = 0 (no restart in that case)
; Waits for boot.sh to complete npm install/build before starting the server.
process_name=%(program_name)s
command=bash /srv/%(ENV_DKZ_DOMAIN)s/.docker/container/scripts/angular.sh
user=root
priority=20
autorestart=unexpected
exitcodes=0
startsecs=0
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
