# .docker/container/etc/supervisor/program/boot-init.conf
[program:boot-init]
; Runs once on first container start (docker run); self-deletes after.
; On subsequent restarts (docker restart) the file is absent — exits cleanly.
process_name=%(program_name)s
command=bash -c "[ -f /etc/dockraft/scripts/boot-init.sh ] && bash /etc/dockraft/scripts/boot-init.sh; exit 0"
user=root
priority=1
autorestart=false
startsecs=0
startretries=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
