# .docker/container/etc/supervisor/program/boot.conf
[program:boot]
; Waits for boot-init.sh to self-delete (signals completion), then runs boot.sh.
process_name=%(program_name)s
command=bash -c "while [ -f /etc/dockraft/scripts/boot-init.sh ]; do sleep 1; done; bash /srv/%(ENV_DKZ_DOMAIN)s/.docker/container/scripts/boot.sh"
user=root
priority=2
autorestart=false
startsecs=0
startretries=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
