# .docker/container/etc/supervisor/program/nginx.conf
[program:nginx]
; nginx is started explicitly by boot.sh (via supervisorctl start nginx) after auth
; configuration is applied. autostart=false prevents nginx from starting before
; boot.sh has replaced #NGINX_BASIC_AUTH_PLACEHOLDER# in the site configs.
; autorestart=true ensures nginx is restarted automatically if it crashes after startup.
process_name=%(program_name)s
command=nginx -g "daemon off;"
user=root
autostart=false
priority=10
autorestart=true
startsecs=1
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
